# prometheus-borg-exporter.timer
[Unit]
Description=Prometheus Borg Exporter Timer

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
